# systemd_units/vizzy-liveness.service
[Unit]
Description=Send liveness ping

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vizzy-liveness.sh
