# systemd_units/vizzy-liveness.timer
[Unit]
Description=Report heartbeat to Vizzybot

[Timer]
OnCalendar=*:0/5
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
